# transport/snmp/target/control/template.service
[Unit]
Description=tb-snmp-transport
After=syslog.target

[Service]
User=thingsboard
ExecStart=/usr/share/tb-snmp-transport/bin/tb-snmp-transport.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
